# deploy/systemd/proxywhirl.service
[Unit]
Description=ProxyWhirl Proxy Rotation Service
After=network.target postgresql.service redis.service
Documentation=https://proxywhirl.readthedocs.io

[Service]
Type=simple
User=proxywhirl
Group=proxywhirl
WorkingDirectory=/opt/proxywhirl

ExecStart=/opt/proxywhirl/.venv/bin/uvicorn proxywhirl.api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=-/etc/default/proxywhirl
Environment="PATH=/opt/proxywhirl/.venv/bin"
Environment="PYTHONUNBUFFERED=1"

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/proxywhirl /var/log/proxywhirl

[Install]
WantedBy=multi-user.target
