# deploy/systemd/proxywhirl.service
[Unit]
Description=ProxyWhirl Proxy Rotation Service
After=network.target postgresql.service redis.service
Documentation=https://proxywhirl.readthedocs.io

[Service]
Type=notify
User=proxywhirl
Group=proxywhirl
WorkingDirectory=/opt/proxywhirl

ExecStart=/opt/proxywhirl/.venv/bin/python -m proxywhirl.api
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/etc/proxywhirl/proxywhirl.conf
Environment="PATH=/opt/proxywhirl/.venv/bin"
Environment="PYTHONUNBUFFERED=1"

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/proxywhirl /var/log/proxywhirl

[Install]
WantedBy=multi-user.target
